# deploy/configs/autostart/connman.desktop
[Desktop Entry]
Encoding=UTF-8
Name=connman-tray
TryExec=connman-gtk --tray
Exec=connman-gtk --tray
Terminal=false
Type=Application