# deploy/configs/autostart/xfce4-clipman-plugin-autostart.desktop
[Desktop Entry]
Hidden=false
RunHook=0
TryExec=mousetweaks --ssc --ssc-time=1.0
Exec=mousetweaks --ssc --ssc-time=1.0