# deploy/configs/autostart/xfce4-clipman-plugin-autostart.desktop
[Desktop Entry]
Hidden=false
RunHook=0
TryExec=xfce4-clipman
Exec=xfce4-clipman
